# devel/etc/catkin/profile.d/50.husky_gazebo.sh
/home/galliumos/kinetic_final/devel/.private/husky_gazebo/etc/catkin/profile.d/50.husky_gazebo.sh